# debian/android-tools-adbd-reboot-watcher-reboot.conf
# Watchdog observing android reboot property.
# This is used by adbd to signal user asked for reboot
# Give android 1 second to complete its reboot before
# taking down also Ubuntu.
# This job watches for a normal reboot request

start on ((android-container sys.powerctl=reboot,) or (android-container sys.powerctl=reboot,*))

task

script
    sleep 1
    /sbin/reboot
end script
