# debian/android-tools-adbd-reboot-watcher-bootloader.conf
# Watchdog observing android reboot property.
# This is used by adbd to signal user asked for reboot
# Give android 1 second to complete its reboot before
# taking down also Ubuntu.
# This job watches for a reboot to bootloader request

start on android-container sys.powerctl=reboot,bootloader

task

script
    sleep 1
    /sbin/reboot -f bootloader
end script
